# bot.conf
[program:AutoFinAdvisorBot]
directory=/home/ubuntu/AutoFinAdvisorBot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/AutoFinAdvisorBot/logfile_err.log
stdout_logfile=/home/ubuntu/AutoFinAdvisorBot/logfile.log